# contrib/init/bitniond.service
# bitniond.service — Systemd unit file for Bitnion daemon
#
# This file enables the bitniond daemon to run as a managed system service
# under systemd. It ensures clean startup, shutdown, logging, and restart behavior.
#
# Place this file in: /etc/systemd/system/bitniond.service
# Then run:
#   sudo systemctl daemon-reexec
#   sudo systemctl enable bitniond
#   sudo systemctl start bitniond
#
# License: MIT
# Adapted from Bitcoin Core systemd integration

[Unit]
Description=Bitnion daemon
After=network.target

[Service]
# Modify this path to point to the installed location of bitniond
ExecStart=/usr/local/bin/bitniond -conf=/etc/bitnion/bitniond.conf -datadir=/var/lib/bitnion

# Restart automatically on failure
Restart=on-failure
PrivateTmp=true
Type=simple
User=bitnion
Group=bitnion
RuntimeDirectory=bitnion
RuntimeDirectoryMode=0710

# System log integration
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bitniond

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_KILL CAP_SYS_CHROOT
ReadWritePaths=/var/lib/bitnion

[Install]
WantedBy=multi-user.target
